# src/cnf.service.in
[Unit]
Description=Update command-not-found catalogs

[Service]
Type=oneshot
ExecStart=/usr/bin/cnf-sync
